# infra/systemd/crontech-api.service
[Unit]
Description=Crontech API (Hono/Bun)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/crontech/apps/api
ExecStart=/root/.bun/bin/bun run dist/index.js
Restart=always
RestartSec=5
EnvironmentFile=/opt/crontech/.env
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
