# services/object-storage/infra/crontech-object-storage.service
# Crontech Self-Hosted Object Storage — BLK-018 v0
#
# Runs the single-node MinIO docker-compose stack defined in
# services/object-storage/docker-compose.yml. The MinIO S3 API listens on
# 127.0.0.1:9000 and the admin console on 127.0.0.1:9001. Public access to
# the S3 API goes through the Bun proxy in services/object-storage/src/index.ts
# (started by a separate systemd unit when wired into production).
#
# Install:
#   cp services/object-storage/infra/crontech-object-storage.service \
#      /etc/systemd/system/crontech-object-storage.service
#   systemctl daemon-reload
#   systemctl enable --now crontech-object-storage
#
# Prerequisites on the host:
#   - Docker Engine + the `docker compose` v2 plugin
#   - /opt/crontech checkout containing services/object-storage/
#   - /opt/crontech/services/object-storage/.env populated with
#     MINIO_ROOT_USER and MINIO_ROOT_PASSWORD
#
# Logs:
#   journalctl -u crontech-object-storage -f

[Unit]
Description=Crontech Self-Hosted Object Storage (MinIO via docker compose)
Documentation=https://docs.crontech.ai/blocks/blk-018
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/crontech/services/object-storage
EnvironmentFile=-/opt/crontech/services/object-storage/.env

ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose pull && /usr/bin/docker compose up -d --remove-orphans

TimeoutStartSec=180
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crontech-object-storage

[Install]
WantedBy=multi-user.target
